# com.cshion.arxivsummaryfeeds.plist.template
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  Template. Do NOT load this file directly (it contains __REPO__ placeholders).
  Generate the real plist with:  ./install_launchd.sh
  which substitutes the repo path and installs it into ~/Library/LaunchAgents.
-->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.cshion.arxivsummaryfeeds</string>

    <!-- A LaunchAgent runs inside the user's GUI login session, so it can
         reach the login Keychain and `claude` stays authenticated (unlike
         cron, which cannot). -->
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>__REPO__/summarize.sh</string>
    </array>

    <!-- summarize.sh uses relative paths (CATEGORY/DATE.md, git add .). -->
    <key>WorkingDirectory</key>
    <string>__REPO__</string>

    <!-- Weekdays (Mon=1 .. Fri=5) at 21:00. launchd fires each dict. -->
    <key>StartCalendarInterval</key>
    <array>
        <dict><key>Weekday</key><integer>1</integer><key>Hour</key><integer>21</integer><key>Minute</key><integer>0</integer></dict>
        <dict><key>Weekday</key><integer>2</integer><key>Hour</key><integer>21</integer><key>Minute</key><integer>0</integer></dict>
        <dict><key>Weekday</key><integer>3</integer><key>Hour</key><integer>21</integer><key>Minute</key><integer>0</integer></dict>
        <dict><key>Weekday</key><integer>4</integer><key>Hour</key><integer>21</integer><key>Minute</key><integer>0</integer></dict>
        <dict><key>Weekday</key><integer>5</integer><key>Hour</key><integer>21</integer><key>Minute</key><integer>0</integer></dict>
    </array>

    <!-- launchd appends to these, so failure history is preserved. -->
    <key>StandardOutPath</key>
    <string>__REPO__/crontab.log</string>
    <key>StandardErrorPath</key>
    <string>__REPO__/crontab.log</string>

    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
